# setrandomhostname/setrandomhostname.service
[Unit]
Description=Set random hostname on each boot
Wants=network-pre.target
Before=network-pre.target

[Service]
ExecStart=/usr/local/bin/setrandomhostname.sh
Type=oneshot

[Install]
WantedBy=multi-user.target
